# apps/sipauthserve.conf
[program:sipauthserve]
command=sudo ./sipauthserve
directory=/OpenBTS/

